# ec2-service/ec2-service.service
[Unit]
Description=EC2 Architecture Generator Service
After=network-online.target ollama.service
Requires=ollama.service

[Service]
Type=simple
User=ec2-service
WorkingDirectory=/opt/ec2-service
Environment="OLLAMA_API_URL=http://localhost:11434"
Environment="OLLAMA_MODEL=deepseek-coder:1.3b-instruct-q4_0"
Environment="PYTHONUNBUFFERED=1"
# Para t2.large: usar --workers 1 (solo 1 worker)
ExecStart=/usr/local/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1 --timeout-keep-alive 1200
# Para t3.xlarge+: usar --workers 2-4
# ExecStart=/usr/local/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Límites de recursos
LimitNOFILE=65536
MemoryMax=8G
CPUQuota=400%

[Install]
WantedBy=multi-user.target

